# config/sardinia_supervisor.conf
[program:sardinia_qtl_browser]
command = /usr/bin/gunicorn server:app -b localhost:5000
directory = /mnt/data/sardinia/
user = ubuntu